# NetNS/scripts/netns-mgr.service
[Unit]
Description=Network Namespace Manager - Restore namespaces on boot
After=network.target
Before=docker.service containerd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/netns-restore.sh restore
ExecStop=/usr/local/bin/netns-restore.sh cleanup
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Run as root (required for netns operations)
User=root
Group=root

# Environment
Environment="HOME=/root"

[Install]
WantedBy=multi-user.target
